# packaging/root/usr/share/faktory/init/faktory.rpm.conf
#
# Upstart script for CentOS 6
#
description "Faktory background jobs"

start on runlevel [2345]
stop on runlevel [016]

umask 0002

# if we crash, restart
respawn
# don't try to restart anymore if we fail 5 times in 5 seconds
respawn limit 5 5

# listen on all interfaces, requires TLS certs in /etc/faktory/tls
# See https://github.com/contribsys/faktory/wiki/Security
# exec /usr/bin/faktory -e production -b :7419

# listen only on localhost
exec /usr/bin/faktory -e production -b localhost:7419
